# os/fs/etc/rc.local
#!/bin/sh
set -eu

sysctl -p /etc/sysctl.conf

hostname oasis &
ip link set lo up &

wait
